# systemd/alphapulse-webapp-ui.service
[Unit]
Description=AlphaPulse Webapp UI (Next.js)
After=network-online.target alphapulse-fastapi.service
Wants=alphapulse-fastapi.service

[Service]
User=alphapulse
WorkingDirectory=/opt/alphapulse/webapp-ui
EnvironmentFile=/opt/alphapulse/.env
ExecStart=/usr/bin/pnpm start -p 3000
Restart=always
RestartSec=5

ProtectSystem=strict
ReadWritePaths=/opt/alphapulse/webapp-ui/.next
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
